# crontab/clear_vote_limit.sh
#!/bin/bash
python manage.py runscript clear_vote_limit
