# crontab/init_vote.sh
#!/bin/bash
python manage.py runscript init_vote